# meta-hometec/recipes-hometec/service/service/statemachine.service
[Unit]
Description=Hometec systemd service.

[Service]
Type=simple
ExecStart=/bin/sh /usr/bin/statemachine.sh

[Install]
WantedBy=multi-user.target