# yocto/meta-hometec/recipes-hometec/service/service/statemachine.service
[Unit]
Description=Example systemd service.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/statemachine.sh

[Install]
WantedBy=multi-user.target